# VPN/crontab
#/var/spool/cron/crontabs/root
@reboot /root/tunnel.sh
@reboot /root/ssocat.sh
#@reboot /usr/sbin/knockd -d
